# shell/sysBuilder/10.8/LaunchAgents/com.apple.CoreRAIDAgent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.CoreRAIDAgent</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.CoreRAID.CoreRAIDAgent</key>
		<true/>
	</dict>
	<key>Program</key>
	<string>/System/Library/PrivateFrameworks/CoreRAID.framework/Versions/A/Resources/CoreRAIDAgent</string>
</dict>
</plist>
